# backend/manas360.service
[Unit]
Description=MANAS360 Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/manas360/backend
ExecStart=/bin/sh -c '/usr/bin/npx prisma migrate deploy && /usr/bin/node dist/server.js'
Restart=always
Environment=NODE_ENV=production
EnvironmentFile=/opt/manas360/backend/.env

[Install]
WantedBy=multi-user.target
